# jupyter.conf
[program:jupyter]
command=bash /local/mnt/workspace/root/jupyter_run.sh
directory=/local/mnt/workspace/root/notebook
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
